# files/system/usr/lib/systemd/system/anubis-setup-user.service
[Unit]
Description=Anubis OS first-boot user setup
After=multi-user.target
ConditionPathExists=!/var/lib/anubis-setup-done

[Service]
Type=oneshot
ExecStart=/usr/share/anubis-os/scripts/setup-ohmybash-user.sh
ExecStartPost=/usr/bin/touch /var/lib/anubis-setup-done
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
